# cron/crontab
*/5 * * * * /usr/bin/docker start cert-scanner >> /var/log/cron.log 2>&1